# lib/svc/manifest/network/network-virtual.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
 Copyright 2017 Erigones, s. r. o.  All rights reserved.
 Use is subject to license terms.

 CDDL HEADER START

 The contents of this file are subject to the terms of the
 Common Development and Distribution License (the "License").
 You may not use this file except in compliance with the License.

 You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE
 or http://www.opensolaris.org/os/licensing.
 See the License for the specific language governing permissions
 and limitations under the License.

 When distributing Covered Code, include this CDDL HEADER in each
 file and include the License file at usr/src/OPENSOLARIS.LICENSE.
 If applicable, add the following below this CDDL HEADER, with the
 fields enclosed by brackets "[]" replaced with your own identifying
 information: Portions Copyright [yyyy] [name of copyright owner]

 CDDL HEADER END

	NOTE:  This service manifest is not editable; its contents will
	be overwritten by package or patch operations, including
	operating system upgrade.  Make customizations in a different
	file.
-->

<service_bundle type='manifest' name='SUNWcsr:network-virtual'>

<service
	name='network/virtual'
	type='service'
	version='1'>

	<dependency name="network-physical"
		grouping="require_all"
		restart_on="none"
		type="service">
		<service_fmri value="svc:/network/physical:default" />
	</dependency>

	<dependency
		name='varpd'
		grouping='require_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/network/varpd' />
	</dependency>

	<dependent
		name='routing-setup'
		grouping='require_all'
		restart_on='none'>
		<service_fmri value='svc:/network/routing-setup' />
	</dependent>

	<dependent
		name='ipfilter'
		grouping='require_all'
		restart_on='none'>
		<service_fmri value='svc:/network/ipfilter' />
	</dependent>

	<dependent
		name='lldpd-server'
		grouping='require_all'
		restart_on='none'>
		<service_fmri value='svc:/network/lldp/server' />
	</dependent>

	<dependent
		name='dns-multicast'
		grouping='require_all'
		restart_on='none'>
		<service_fmri value='svc:/network/dns/multicast' />
	</dependent>

	<dependent
		name='ibd-post-upgrade'
		grouping='require_all'
		restart_on='none'>
		<service_fmri value='svc:/network/ibd-post-upgrade' />
	</dependent>

	<dependent
		name='location'
		grouping='require_all'
		restart_on='none'>
		<service_fmri value='svc:/network/location' />
	</dependent>

	<dependent
		name='iptun'
		grouping='require_all'
		restart_on='none'>
		<service_fmri value='svc:/network/iptun' />
	</dependent>

	<dependent
		name='ssh'
		grouping='require_all'
		restart_on='none'>
		<service_fmri value='svc:/network/ssh' />
	</dependent>

	<dependent
		name='vrrp'
		grouping='require_all'
		restart_on='none'>
		<service_fmri value='svc:/network/vrrp' />
	</dependent>

	<dependent
		name='milestone-network'
		grouping='optional_all'
		restart_on='restart'>
		<service_fmri value='svc:/milestone/network' />
	</dependent>

	<dependent
		name='smartdc-init'
		grouping='require_all'
		restart_on='none'>
		<service_fmri value='svc:/system/smartdc/init' />
	</dependent>

	<dependent
		name='identity'
		grouping='optional_all'
		restart_on='none'>
		<service_fmri value='svc:/system/identity' />
	</dependent>

	<instance name='default' enabled='true'>

	<exec_method
		type='method'
		name='start'
		exec='/lib/svc/method/net-virtual start'
		timeout_seconds='300' />

	<exec_method
		type='method'
		name='stop'
		exec=':true'
		timeout_seconds='3' />

	<exec_method
		type='method'
		name='refresh'
		exec='/lib/svc/method/net-virtual refresh'
		timeout_seconds='300' />

	<property_group name='startd' type='framework'>
		<propval name='duration' type='astring' value='transient' />
		<propval name='critical_failure_count' type='integer' value='1'/>
	</property_group>

	<template>
		<common_name>
			<loctext xml:lang='C'>
				virtual network interfaces
			</loctext>
		</common_name>
		<documentation>
			<manpage title='ifconfig' section='1M'
				manpath='/usr/share/man' />
		</documentation>
	</template>

	</instance>

	<stability value='Unstable' />

</service>

</service_bundle>

